# data/conf/supervisor/workers.conf
[group:workers]
programs=consumerFibbonacci,consumerSimply

[program:consumerFibbonacci]
command=php index.php -t consumer -n fibonacci
directory=/home/testuser
numprocs=1
autorestart=true
startsecs=5
user=testuser

[program:consumerSimply]
command=php index.php -t consumer -n simply
directory=/home/testuser
numprocs=1
autorestart=true
startsecs=5
user=testuser

[group:generators]
programs=producerFibbonacci,producerSimply

[program:producerFibbonacci]
command=php index.php -t producer -n fibonacci --amount 2000 --usleep 100000
directory=/home/testuser
numprocs=1
autorestart=false
startsecs=5
user=testuser

[program:producerSimply]
command=php index.php -t producer -n simply --amount 5000 --usleep 150000
directory=/home/testuser
numprocs=1
autorestart=false
startsecs=5
user=testuser